# roles/gohub/templates/gohub.conf
start on net-device-up
exec /usr/bin/gohub --log=/var/log/gohub.log --port=7654 --config=/etc/gohub.json --secret={{ github_secret }}
